# etc/systemd/system/nginx.service
[Unit]
Description=Nginx web server

[Service]
ExecStart=/srv/nginx/sbin/nginx
Type=forking
PIDFile=/srv/nginx/logs/nginx.pid
